# rootfs/etc/services.d/saslauthd/run
#!/usr/bin/execlineb -P

saslauthd -d -a pam