# rootfs/etc/services.d/svnserve/run
#!/usr/bin/execlineb -P

svnserve -d --foreground --config-file=/etc/svnserve.conf -r /data/repos